# crontab/crontab
* * * * * BASH_ENV=/app/env.sh /app/task.sh >> /app/task.log 2>&1
