# crontab/env.sh
#!/bin/bash

# export OPENAI_API_KEY=$OPENAI_API_KEY
# export LINKEDIN_ORG_ID=$LINKEDIN_ORG_ID
# export LINKEDIN_ACCESS_TOKEN=$LINKEDIN_ACCESS_TOKEN